# com.tennis.wakeup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tennis.wakeup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/willwallwan/Documents/GitHub/Octogon/ensure_wake_and_run.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardErrorPath</key>
    <string>/Users/willwallwan/Documents/GitHub/Octogon/launch_error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/willwallwan/Documents/GitHub/Octogon/launch_output.log</string>
</dict>
</plist> 